# deploy/systemd/td-bump-overdue.service
[Unit]
Description=Update overdue Toodledo task dates to today

[Service]
Type=oneshot
WorkingDirectory=/home/crog/Projects/TD-Toodle-App
ExecStart=/home/crog/Projects/TD-Toodle-App/.venv/bin/td bump-overdue --apply
